# system-patches/disable-wifi-aspm.service
[Unit]
Description=Disable ASPM for Intel Wi-Fi (Fix Instability)
After=multi-user.target suspend.target

[Service]
Type=oneshot
# This command reads the current state, forcibly turns off the ASPM bits (last 2 bits), and saves it.
ExecStart=/bin/bash -c 'val=$(setpci -s 00:14.3 CAP_EXP+10.w); new_val=$(printf "%%04x" $((0x$val & 0xFFFC))); setpci -s 00:14.3 CAP_EXP+10.w=$new_val'

[Install]
WantedBy=multi-user.target suspend.target
